# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:api]
command=uvicorn app.main:app --host 0.0.0.0 --port %(ENV_PORT)s --log-level info
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/api.err.log
stdout_logfile=/var/log/supervisor/api.out.log

[program:worker]
command=celery -A app.workers.celery_app worker --queues=documents,default --concurrency=2 --loglevel=info --without-heartbeat --without-gossip --without-mingle
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/worker.err.log
stdout_logfile=/var/log/supervisor/worker.out.log
stopwaitsecs=10
